# common/.config/autostart/run-on-startup.desktop
[Desktop Entry]
Name=Run on startup
Type=Application
Exec=/usr/bin/bash /home/bharatkalluri/.run_on_startup
Terminal=false